# systemd/minio-compose.service
[Unit]
Description=MinIO stack via Docker Compose
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/minio-setup
ExecStart=/usr/bin/docker compose up -d minio
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
